# restic.timer
[Unit]
Description=Restic Backup Timer

[Timer]
WakeSystem=false
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
